# deploy/systemd/smartdisplay-x11.service
[Unit]
Description=X11 Display Server for Smart Display
After=systemd-user-sessions.service

[Service]
ExecStart=/usr/bin/Xorg :0 vt1 -nolisten tcp -s 0 dpms
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
